# deploy/board.service
[Unit]
Description=Board FastAPI backend
After=network.target postgresql.service signal-cli.service

[Service]
Type=simple
User=kyle46220
WorkingDirectory=/opt/backend
EnvironmentFile=/opt/backend/.env
ExecStart=/opt/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips='*'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
